# tools/garbage-collector.cron
0 * * * * passhport /home/passhport/passhport/tools/garbage-collector.sh
